# root/var/spool/cron/crontabs/rancid
# m h dom mon dow usercommand
*/5 * * * * rancid-run
